# internal/svcinstall/templates/systemd.service.tmpl
[Unit]
Description=Tailscale Metrics Discovery — {{.Role}}
After=network.target
Wants=tailscaled.service

[Service]
Type=simple
ExecStart={{.BinaryPath}} {{.Role}} daemon --config {{.ConfigFile}}
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
